# proj1/src/main/resources/Dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:node-app]
command=/app/startup.sh
user=node

[program:log-rotation]
command=/bin/sh -c "/app/logrotate.sh"
user=node